# scalets-backup.timer
[Unit]
Description=backup scalets in vscale.io

[Timer]
OnCalendar=*-*-* *:00:00
Unit=scalets-backup.service

[Install]
WantedBy=multi-user.target